# templates/plotter-certrenew.template.timer
[Unit]
Description=Renew cert once a day

[Timer]
OnCalendar=daily

[X-Fleet]
MachineMetadata=n=XX_NODENAME
